# scripts/crontab
# Crontab for automatic backups
# Runs backup script every day at 2 AM

0 2 * * * /backup.sh >> /var/log/backup.log 2>&1

# Weekly backup cleanup (remove backups older than 30 days)
0 3 * * 0 find /backups -name "backup_*.sql.gz" -mtime +30 -delete

# Monthly database optimization (vacuum and analyze)
0 4 1 * * PGPASSWORD=$POSTGRES_PASSWORD psql -h db -U $POSTGRES_USER -d $POSTGRES_DB -c "VACUUM ANALYZE;" >> /var/log/maintenance.log 2>&1